# kali-config/variant-mate/includes.chroot/etc/systemd/system/xinit-login.service
[Unit]
After=systemd-user-sessions.service
 
[Service]
ExecStart=/bin/su root -l -c /usr/bin/startx -- VT01
 
[Install]
WantedBy=multi-user.target
